# systemd/backup.timer
[Unit]
Description=Daily backup timer for amstore

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
